# cookbooks/webobjects/templates/default/wo-webobjects.initd.erb
#!/bin/bash

# chkconfig: - 90 20
# description: Provides WebObjects services

<% if %w{ redhat centos scientific fedora }.include?(node[:platform]) -%>
# Source function library.
. /etc/rc.d/init.d/functions
<% end -%>

NEXT_ROOT="<%= node[:webobjects][:webobjects_WOLocalRootDirectory_dir] %>"
export NEXT_ROOT

USER=<%= node[:webobjects][:webobjects_user] %>

# See how we were called.
case "$1" in
	start)
		echo -n "Starting wotaskd and Monitor: "
		su $USER -c "<%= node[:webobjects][:webobjects_WODeployment_dir] %>/wotaskd.woa/wotaskd -WOPort <%= node[:webobjects][:webobjects_wotaskd_port] %> &"
		su $USER -c "<%= node[:webobjects][:webobjects_WODeployment_dir] %>/JavaMonitor.woa/JavaMonitor -WOPort <%= node[:webobjects][:webobjects_JavaMonitor_port] %> &"
		echo
		;;
	stop)
		echo -n "Shutting down wotaskd and Monitor: "
		WOTASKD_PID=`ps aux | awk '/WOPort <%= node[:webobjects][:webobjects_wotaskd_port] %>/ && !/awk/ {print $2}'`
		kill $WOTASKD_PID
		MONITOR_PID=`ps aux | awk '/WOPort <%= node[:webobjects][:webobjects_JavaMonitor_port] %>/ && !/awk/ {print $2}'`
		kill $MONITOR_PID
		echo
		;;
	restart)
		$0 stop
		$0 start
		;;
	*)
		echo -n "Usage: $0 {start|stop|restart}"
		exit 1
esac

if [ $# -gt 1 ]; then
	shift
	$0 $*
fi

exit 0